# system-files/etc/systemd/system/greetd.service.d/after-zena-setup.conf
[Unit]
Wants=zena-setup.service
After=zena-setup.service
