# system-files/base/etc/profile.d/dms-setup.sh
#!/bin/sh

DMS="$HOME/.config/niri/dms"

mkdir -p "$DMS"

touch "$DMS/colors.kdl"
touch "$DMS/wpblur.kdl"
touch "$DMS/cursor.kdl"
touch "$DMS/outputs.kdl"
touch "$DMS/layout.kdl"
touch "$DMS/alttab.kdl"
touch "$DMS/binds.kdl"
